# root/etc/init.d/blockclient
#!/bin/sh /etc/rc.common

START=99
USE_PROCD=1

start_service() {
    procd_open_instance
    procd_set_param command /bin/true
    procd_close_instance
    apply_rules
}

reload_service() {
    apply_rules
}

apply_rules() {
    iptables -F BLOCKCLIENT 2>/dev/null || iptables -N BLOCKCLIENT
    iptables -D FORWARD -j BLOCKCLIENT 2>/dev/null
    iptables -I FORWARD -j BLOCKCLIENT

    uci -q show blockclient | awk -F= '/=rule$/ {print $1}' | while read section; do
        ip=$(uci -q get "$section.ip")
        mac=$(uci -q get "$section.mac")

        [ -n "$ip" ] && iptables -A BLOCKCLIENT -s "$ip" -j DROP
        [ -n "$mac" ] && iptables -A BLOCKCLIENT -m mac --mac-source "$mac" -j DROP
    done
}
